# presence/service_files/json/phosphor-fan-presence-tach@.service
[Unit]
Description=Phosphor Fan Presence Tach Daemon
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Restart=on-failure
ExecStart=@BINDIR@/phosphor-fan-presence-tach
