# presence/service_files/json/phosphor-fan-presence-tach@.service.in
[Unit]
Description=Phosphor Fan Presence Tach Daemon

[Service]
Restart=on-failure
ExecStart=@BINDIR@/phosphor-fan-presence-tach
