# monitor/service_files/json/phosphor-fan-monitor@.service.in
[Unit]
Description=Phosphor Fan Monitor Daemon
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Restart=on-failure
ExecStart=@BINDIR@/phosphor-fan-monitor
