# monitor/service_files/json/phosphor-fan-monitor@.service.in
[Unit]
Description=Phosphor Fan Monitor Daemon

[Service]
Restart=on-failure
ExecStart=@BINDIR@/phosphor-fan-monitor
