# control/service_files/json/phosphor-fan-control@.service
[Unit]
Description=Phosphor Fan Control Daemon

[Service]
Restart=on-failure
ExecStart=/usr/bin/phosphor-fan-control
