# control/service_files/json/phosphor-fan-control@.service.in
[Unit]
Description=Phosphor Fan Control Daemon

[Service]
Restart=on-failure
ExecStart=@BINDIR@/phosphor-fan-control
